# provisioning/files/selenium.conf
# Starts up selenium for browser testing
description     "regular background program processing daemon"

start on runlevel [2345]
stop on runlevel [!2345]

expect fork
respawn

exec /bin/selenium-start.sh